# dev-python/carbon/files/carbon-cache
#!/sbin/runscript

depend() {
	use network
}

start() {
	/opt/graphite/bin/carbon-cache.py start
	eend $?
}

stop() {
	/opt/graphite/bin/carbon-cache.py stop
	eend $?
}

status() {
	/opt/graphite/bin/carbon-cache.py status
	eend $?
}

